# roles/base/files/rc.local
#!/bin/sh
cgconfigparser -l /etc/cgconfig.conf
cgrulesengd